# scripts/sicos-usb-monitor.service
[Unit]
Description=SICOS USB Monitor Service
After=network.target

[Service]
ExecStart=/full/path/to/scripts/usb_monitor.sh
Restart=always
User=your-username

[Install]
WantedBy=multi-user.target
